# supervisord.conf
[supervisord]
nodaemon = true
[program:3proxy]
command  = 3proxy /3proxy.cfg

[program:socat]
command  = socat UNIX-LISTEN:/vpn/3proxy.sock,fork TCP:localhost:3128